# config/ai.openclaw.daily-news.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.daily-news</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/leonard/.openclaw/workspace/scripts/daily_news_automation.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/leonard/.openclaw/workspace</string>
    <key>StandardOutPath</key>
    <string>/Users/leonard/.openclaw/workspace/logs/daily-news.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/leonard/.openclaw/workspace/logs/daily-news-error.log</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
</dict>
</plist>
